# backend/code/etc/init/cronrat.conf

#!upstart
description "cronrat cluster server"
author      "Andrew Yasinsky"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/root"

    exec `/var/cronrat_bin/cronrat-server -Configfile="/var/cronrat_bin/cronrat-server.json" >> /var/log$
end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/cronrat/cronrat.log
    echo `status cronrat | egrep -oi '([0-9]+)$' | head -n1` > /var/run/cronrat.pid
end script

pre-stop script
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/cronrat/cronrat.log
    rm /var/run/cronrat.pid
end script